# docker/cal2wl/cron/crontab
30 * * * * /cron/dicom-labbooking2worklist.sh
00 15 * * * /cron/dicom-clean-worklist.sh
